# ai-school-web.service
[Unit]
Description=AI School Web Application
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/basel/Projects/AI-school/web
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /home/basel/Projects/AI-school/web/.next/standalone/server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-school-web

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/basel/Projects/AI-school/web

[Install]
WantedBy=multi-user.target
